# ansible/roles/visibike.jsondumpster/templates/autopostgresqlbackup.service
[Unit]
Description=Backups done with autopostgresqlbackup

[Service]
Type=oneshot
ExecStart=/usr/sbin/autopostgresqlbackup